# ops/systemd/collab-tracker.service
[Unit]
Description=TaskTrack - Collaborative Task Tracker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/rtoony/projects/collab-tracker
ExecStart=/home/rtoony/miniconda3/bin/gunicorn --bind 0.0.0.0:5050 --workers 2 --access-logfile - app:app
Restart=on-failure
RestartSec=5
Environment=PATH=/home/rtoony/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
