# ops/systemd/collab-tracker.service
[Unit]
Description=TaskTrack - Collaborative Task Tracker
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/rtoony/projects/collab-tracker
ExecStartPre=/home/rtoony/bin/nexus-svc-inject collab-tracker "TaskTrack" "Infrastructure"
ExecStart=/home/rtoony/miniconda3/bin/gunicorn -c /home/rtoony/projects/collab-tracker/gunicorn.conf.py wsgi:app
ExecStopPost=/bin/rm -f /dev/shm/nexus-env-collab-tracker
EnvironmentFile=-/dev/shm/nexus-env-collab-tracker
Restart=on-failure
RestartSec=5
Environment=PATH=/home/rtoony/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
